# os/runit/sv/snooze/run
#!/bin/sh
exec 2>&1
exec snooze -H0-23 -M0 -- run-parts /etc/cron.hourly
