# os/runit/sv/agetty-tty2/run
#!/bin/sh
exec 2>&1
exec agetty tty2 38400 linux
