# os/runit/sv/udevd/finish
#!/bin/sh
exec 2>&1
udevadm control --exit 2>/dev/null
exit 0
