# os/runit/sv/agetty-tty1/run
#!/bin/sh
exec 2>&1
exec agetty --noclear tty1 38400 linux
