# os/runit/sv/dhcpcd/run
#!/bin/sh
exec 2>&1
[ -d /run/dhcpcd ] || mkdir -p /run/dhcpcd
exec dhcpcd -B -M
